# setup/nginx_config_example/uwsgi-robot.service
[Unit]
Description=uWSGI service for robot
After=syslog.target

[Service]
ExecStartPre=/bin/mkdir -p -m 775 /var/run/uwsgi
ExecStart=/usr/local/bin/uwsgi --ini /etc/uwsgi/robot.ini
RuntimeDirectory=/var/run/uwsgi
Restart=always
KillSignal=SIGTERM
StandardError=syslog
Type=forking
NotifyAccess=all

[Install]
WantedBy=multi-user.target
